# etc/profile.d/node_app_env.sh
export token=token-in-etc-profile.d-env.sh
